# roles/yh-run-clean/templates/yh-run-clean.service.j2
[Unit]
Description={{ yhrc_service_name }}
After=network.target

[Service]
Type=forking
ExecStart={{ yhrc_service_ctl_script }} start
PrivateTmp=true

[Install]
WantedBy=multi-user.target
